# notes/Linux/fhs/etc/inittab
$ /etc/inittab

This it the toplevel configuration file for (SysV style) → Linux/init.

sa:
  → Linux/fhs/etc
