# notes/Linux/fhs/etc/inittab
$ /etc/inittab

This it the toplevel configuration file for (→ Linux/init/System-V style) → Linux/init.

sa:
  → Linux/fhs/etc
