# notes/Linux/fhs/etc/inittab
$ /etc/inittab

This it the toplevel configuration file for (→ Linux/init/System-V style) → Linux/init. It describes how the init process should set up the
system in a certain *run-level*.

code {
#
#    The initdefault entry determines the initial
#    runlevel of the system.
#
id:3:initdefault:
code }



sa:
  → Linux/fhs/etc
