# notes/Linux/fhs/etc/inittab
$ /etc/inittab

This it the toplevel configuration file for → Linux/init.
